# sinks/cloudflow-sink-splunk/systemd/cloudflow-sink-splunk.service
[Unit]
Description=CloudFlow Splunk sink
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Must define SPLUNK_HEC_TOKEN (or whatever env var splunk.hec_token_env in
# the YAML config names) -- decision D6: the HEC token is env-only, never in
# YAML. File is not committed; deploy it out of band, e.g.:
#   SPLUNK_HEC_TOKEN=...
EnvironmentFile=-/etc/cloudflow/splunk-sink.env
# WP-17: ExecStart is the compiled C binary (was the Python console script).
ExecStart=/usr/local/bin/cloudflow-sink-splunk -c /etc/cloudflow/splunk-sink.yaml
# SIGTERM triggers a flush-once-then-exit; give the in-flight batch time.
KillSignal=SIGTERM
TimeoutStopSec=15
Restart=on-failure
RestartSec=5
User=cloudflow
Group=cloudflow
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
